# ImgProc.service
[Unit]
Description=Image Processing service

[Service]
Type=oneshot
ExecStart=/home/root/IntUnderflow

[Install]
WantedBy=multi-user.target
